# systemd/media-dir.service
[Unit]
Description=Make media dir

[Service]
User=root
Type=simple
ExecStart=/root/.scripts/media-dir.sh
